# ops/systemd/trading-ops-db-backup.timer
[Unit]
Description=Trading Ops - Daily Postgres backup timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target

